# odi_production/services/odi-cortex-query.service
[Unit]
Description=ODI Cortex Query - API RAG Multi-Lóbulo
After=network.target
Documentation=https://github.com/juandavidjd/extrac

[Service]
Type=simple
User=root
WorkingDirectory=/opt/odi
ExecStart=/usr/bin/python3 -m uvicorn core.odi_cortex_query:app --host 0.0.0.0 --port 8803
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Logs
StandardOutput=append:/opt/odi/logs/cortex_query_stdout.log
StandardError=append:/opt/odi/logs/cortex_query_stderr.log

[Install]
WantedBy=multi-user.target
